# setup/sysmon-oled.service
[Unit]
Description=System monitoring I2C OLED display daemon

[Service]
Type=simple
ExecStart=/bin/bash /opt/sysmon-oled/launch.sh

[Install]
WantedBy=multi-user.target
